# mobile-detector.service
[Unit]
Description=Mobile Data Detector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/mobile-detector-venv/bin/python /usr/local/bin/mobile_detector start
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
User=root
Group=root

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/lib/mobile_data_monitor /tmp

[Install]
WantedBy=multi-user.target
